# telegram-bot-api.service
[Unit]
Description=Telegram Bot API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt
Environment=TELEGRAM_API_ID=93618
Environment=TELEGRAM_API_HASH=6ac1a760c900b2c0bdd054f15e3522d5
ExecStart=/opt/telegram-bot-api --local --http-port=8081 --dir=/opt/telegram-bot-api-data --temp-dir=/opt/telegram-bot-api-data/tmp
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/telegram-bot-api-data

[Install]
WantedBy=multi-user.target
