# telegram-bot-api.service
[Unit]
Description=Telegram Bot API Server (Official TDLib Build)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/telegram-bot-api-data
Environment=TELEGRAM_API_ID=YOUR_API_ID_HERE
Environment=TELEGRAM_API_HASH=YOUR_API_HASH_HERE
ExecStart=/usr/local/bin/telegram-bot-api --local --http-port=8081 --dir=/opt/telegram-bot-api-data --temp-dir=/opt/telegram-bot-api-data/tmp
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/telegram-bot-api-data

[Install]
WantedBy=multi-user.target
